# root/etc/init.d/luci-app-tailscale-ng
#!/bin/sh /etc/rc.common

# Copyright (C) 2026 vad-b
# SPDX-License-Identifier: GPL-3.0-only

# This init script manages 'tailscale up' configuration options.
# It does NOT manage the tailscaled daemon itself - that's handled by
# the standard /etc/init.d/tailscale from the tailscale package.
#
# This script reads UCI options from /etc/config/luci-app-tailscale-ng
# and applies them via 'tailscale up --reset'.

USE_PROCD=1
START=95

service_triggers() {
	procd_add_reload_trigger "tailscale"
	procd_add_reload_trigger "luci-app-tailscale-ng"
}

start_service() {
	logger -t "luci-app-tailscale-ng" "Triggered: start_service"

	local enabled enabled_status ts_service_status
	local login_server authkey
	local accept_dns accept_routes advertise_exit_node
	local exit_node advertise_routes flags

	config_load "luci-app-tailscale-ng"
	config_get_bool enabled "settings" "enabled" 0


	# Handle enabled/disabled state
	[ "$enabled" -eq 1 ] && enabled_status="true" || enabled_status="false"
	logger -t "luci-app-tailscale-ng" "Tailscale service enabled: [$enabled_status]"
	if [ "$enabled" -eq 0 ]; then
		logger -t "luci-app-tailscale-ng" "Tailscale service disabled, stopping tailscale"
		/etc/init.d/tailscale stop
		return 0
	fi

	# Start/reload tailscale service if status is not "running"
	ts_service_status="$(/sbin/service tailscale status)"
	logger -t "luci-app-tailscale-ng" "Tailscale service status: [$ts_service_status]"
	if [ "$ts_service_status" != "running" ]; then
		logger -t "luci-app-tailscale-ng" "Tailscale service is not running, reloading tailscale"
		/etc/init.d/tailscale reload
	fi

	# Read configuration options
	config_get login_server "settings" "login_server" ""
	config_get authkey "settings" "authkey" ""
	config_get_bool accept_dns "settings" "accept_dns" 1
	config_get_bool accept_routes "settings" "accept_routes" 0
	config_get_bool advertise_exit_node "settings" "advertise_exit_node" 0
	config_get exit_node "settings" "exit_node" ""
	config_get advertise_routes "settings" "advertise_routes" ""
	config_get flags "settings" "flags" ""

	#
	# Build and execute tailscale up via procd-managed instance
	#
	local routes_csv
	set -- /usr/sbin/tailscale up --reset

	[ -n "$login_server" ] && set -- "$@" "--login-server=$login_server"
	[ -n "$authkey" ] && set -- "$@" "--authkey=$authkey"
	[ "$accept_dns" = "0" ] && set -- "$@" "--accept-dns=false"
	[ "$accept_routes" = "1" ] && set -- "$@" "--accept-routes=true"
	[ "$advertise_exit_node" = "1" ] && set -- "$@" "--advertise-exit-node"

	[ -n "$exit_node" ] && {
		set -- "$@" "--exit-node=$exit_node"
		set -- "$@" "--exit-node-allow-lan-access=true"
	}

	[ -n "$advertise_routes" ] && {
		routes_csv=$(echo "$advertise_routes" | tr ' ' ',')
		set -- "$@" "--advertise-routes=$routes_csv"
	}

	[ -n "$flags" ] && set -- "$@" $flags

	logger -t "luci-app-tailscale-ng" "Executing: $*"
	procd_open_instance up
	procd_set_param command "$@"
	procd_set_param stdout 1
	procd_set_param stderr 1
	procd_close_instance
	return 0
}

stop_service() {
	logger -t "luci-app-tailscale-ng" "Stopping tailscale"
	/etc/init.d/tailscale stop
}
